# system/overlays/tinc/etc/systemd/system/symrouter.service
# -*- conf -*-

[Unit]
Description=symrouter
After=network-pre.target
Before=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/tincd -n symrouter

[Install]
WantedBy=multi-user.target
