# system/overlays/network/etc/systemd/system/hostapd@.service.d/override.conf
[Unit]
After=sys-subsystem-net-devices-%i.device
